# docker/certupdate.crontab
_MICRON_SYSLOG_FACILITY = cron
@reboot		root	/proxycam/certupdate
0 * * * *	root	/proxycam/certupdate
